# wifi-usb/apkovl/etc/init.d/gadget-ha-rclone
#!/sbin/openrc-run

name="gadget-ha-rclone"
description="USB Gadget Home Assistant + Rclone Sync"

command="/usr/bin/gadget-ha-rclone"
command_background="yes"
pidfile="/run/${name}.pid"
output_logger="logger -t ${name} -p daemon.info"
error_logger="logger -t ${name} -p daemon.err"

export CONFIG_URL="@CONFIG_URL@"
export RCLONE_CONFIG="@RCLONE_CONFIG@"

depend() {
    need net localmount
    after usb-gadget
}
