# rpi-authbox/rootfs-overlay/etc/systemd/system/serial-getty@.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --keep-baud 115200 --autologin root --noclear - $TERM
