# rpi-authbox/rootfs-overlay/etc/systemd/system/systemd-journal-upload.service.d/install.conf
[Service]
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
